# package/gluon-wan-dnsmasq/files/etc/hotplug.d/iface/50-gluon-wan-dnsmasq
if grep -qxFe "$INTERFACE" /lib/gluon/wan-dnsmasq/interface.d/*; then
	/lib/gluon/wan-dnsmasq/update.lua
fi
