# package/gluon-respondd/files/etc/hotplug.d/iface/10-gluon-respondd
#!/bin/sh

. /usr/share/libubox/jshn.sh
. /lib/functions/service.sh

DEVLIST="$(uci_get "respondd.@respondd[0].iface_list_file")"

ifname_to_dev () {
	json_load "$(ubus call network.interface.$1 status)"
	json_get_var dev device

	echo "$dev"
}

case "$ACTION" in
	ifdown)
		DEVICE="$(ifname_to_dev "$INTERFACE")"
		sed "/^$DEVICE$/d" $DEVLIST > $DEVLIST.new
		mv $DEVLIST.new $DEVLIST
		;;
	ifup)
		DEVICE="$(ifname_to_dev "$INTERFACE")"
		MESH="$(cat "/sys/class/net/$DEVICE/batman_adv/mesh_iface" 2>/dev/null)"

		[ "$MESH" = "bat0" -o "$INTERFACE" = "client" ] || exit 0

		{
			cat $DEVLIST 2>/dev/null
			echo $DEVICE
		} | sort -u > $DEVLIST.new
		mv $DEVLIST.new $DEVLIST

		/etc/init.d/respondd reload &

		;;
esac
