# package/gluon-respondd/files/etc/hotplug.d/iface/10-gluon-respondd
#!/bin/sh

[ "$ACTION" = "ifupdate" ] && exit

exec /etc/init.d/gluon-respondd reload
