# package/gluon-respondd/files/etc/hotplug.d/iface/10-gluon-respondd
#!/bin/sh

. /lib/functions/service.sh

case "$ACTION" in
	ifdown) ;;
	ifup)
		/etc/init.d/gluon-respondd restart_if_running &
		;;
esac
