# package/gluon-radvd/files/etc/init.d/gluon-radvd
#!/bin/sh /etc/rc.common

START=50

SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1


start() {
	service_start /usr/sbin/uradvd $(/lib/gluon/radvd/arguments)
}

stop() {
	service_stop /usr/sbin/uradvd
}
