# package/gluon-radvd/files/etc/init.d/gluon-radvd
#!/bin/sh /etc/rc.common

START=50

SERVICE_USE_PID=1
SERVICE_NAME=gluon-radvd
SERVICE_PID_DIR=/var/run/gluon-radvd
SERVICE_PID_FILE="$SERVICE_PID_DIR"/gluon-radvd.pid

radvd_conf=/var/gluon/radvd/radvd.conf

start() {
	mkdir -p "$SERVICE_PID_DIR"
	chown gluon-radvd "$SERVICE_PID_DIR"

	mkdir -p "$(dirname "$radvd_conf")"
	/lib/gluon/radvd/generate_config > "$radvd_conf"

	service_start /usr/sbin/radvd -C "$radvd_conf" -m stderr_syslog -u gluon-radvd -p "$SERVICE_PID_FILE"
}

stop() {
	service_stop /usr/sbin/radvd
	rm "$radvd_conf"
}
