# package/gluon-luci-theme/files/etc/uci-defaults/luci-theme-gluon
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Gluon=/luci-static/gluon
        commit luci
EOF
	
