# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S12rpcd
#!/bin/sh /etc/rc.common

if [ -x /etc/init.d/rpcd ]; then
	. /etc/init.d/rpcd
fi
