# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S13haveged
#!/bin/sh /etc/rc.common

if [ -x /etc/init.d/haveged ] && /etc/init.d/haveged enabled; then
	. /etc/init.d/haveged
fi
