# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S13haveged
#!/bin/sh /etc/rc.common

if /etc/init.d/haveged enabled; then
	. /etc/init.d/haveged
fi
