# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S13urngd
#!/bin/sh /etc/rc.common

# shellcheck disable=SC1091

if [ -x /etc/init.d/urngd ] && /etc/init.d/urngd enabled; then
	. /etc/init.d/urngd
fi
