# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S50telnet
#!/bin/sh /etc/rc.common

START=50

USE_PROCD=1
PROG=/usr/sbin/telnetd

start_service() {
	procd_open_instance
	procd_set_param command "$PROG" -F -l /lib/gluon/setup-mode/ash-login
	procd_close_instance
}
