# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S96led
#!/bin/sh /etc/rc.common

START=96

start() {
	/etc/init.d/led start

	. /etc/diag.sh
	get_status_led
	status_led_set_timer 1000 300
}
