# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S96led
#!/bin/sh /etc/rc.common

START=96

start() {
	/etc/init.d/led start
	/lib/gluon/setup-mode/led.sh running
}
