# package/gluon-mesh-batman-adv/files/etc/profile.d/001-gluon-mesh-batman-adv.sh
#!/bin/sh

export PATH="/usr/local/sbin:$PATH"
