# package/gluon-setup-mode/files/lib/gluon/setup-mode/init.d/bootcount
#!/bin/sh /etc/rc.common

# shellcheck disable=SC1091

if [ -x /etc/init.d/bootcount ] && /etc/init.d/bootcount enabled; then
	. /etc/init.d/bootcount
fi
