# package/gluon-core/files/etc/uci-defaults/zzz-gluon-upgrade
#!/bin/sh

for script in /lib/gluon/upgrade/*; do
	"$script"
done
