# package/gluon-mesh-vpn-fastd/files/etc/hotplug.d/iface/50-gluon-mesh-vpn-fastd
[ "$INTERFACE" = 'wan' ] || exit 0

/lib/gluon/mesh-vpn-fastd/update-dnsmasq
