# package/gluon-mesh-babel/files/etc/hotplug.d/iface/90-gluon-babel-brclient-configure-snooping
#!/bin/sh

case "$ACTION"
in
	ifup)
	echo 1 > /sys/devices/virtual/net/br-client/bridge/multicast_snooping
	echo 2 > /sys/devices/virtual/net/br-client/bridge/multicast_router
	;;
esac
